# home/.config/systemd/user/fcitx5.service
[Unit]
Description="fcitx5 service"

[Service]
ExecStart=/usr/bin/fcitx5 -D

[Install]
WantedBy=default.target
